# versions/7.2/runtime/config/supervisor/cron.conf
[program:cron]
command=/usr/sbin/cron -f -L 7
